# infrastructure/conf/reverb.conf
[program:laravel-reverb]
command=php /home/site/wwwroot/api/artisan reverb:start
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/laravel-reverb.log
